# dist/P4wnP1-kali.service
[Unit]
Description=P4wnP1 Startup Service
#After=systemd-modules-load.service
After=local-fs.target
DefaultDependencies=no
Before=sysinit.target

[Service]
Type=simple
#Type=forking
#RemainAfterExit=yes
ExecStartPre=/usr/local/bin/monassure
ExecStart=/usr/local/bin/P4wnP1_service
#StandardOutput=journal+console
#StandardError=journal+console
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
#WantedBy=sysinit.target
